# scripts/com.timestripe.mailbox-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.timestripe.mailbox-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/vladislavrumancev/Desktop/Cursor/Timestripe/.venv/bin/python</string>
        <string>/Users/vladislavrumancev/Desktop/Cursor/Timestripe/reply_detector.py</string>
        <string>--apply</string>
        <string>--since</string>
        <string>2</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/vladislavrumancev/Desktop/Cursor/Timestripe</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/vladislavrumancev/Desktop/Cursor/Timestripe/output/mailbox_sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/vladislavrumancev/Desktop/Cursor/Timestripe/output/mailbox_sync.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
